# scripts/sync-schedule.service.default
[Unit]
Description=synchronise conferenceTimer every minute

[Service]
User=pi
ExecStart=/path/to/conferenceTimer/scripts/sync.sh
