# scripts/sync-schedule.service
[Unit]
Description=synchronise conferenceTimer every minute

[Service]
User=pi
ExecStart=$HOME/conferenceTimer/scripts/sync.sh
